# etc/systemd/system/shared-config-init.service.d/before-shutdown.conf
[Unit]
DefaultDependencies=no
Before=shutdown.target
