# etc/systemd/system/shared-config-init.service.d/reload-systemd.conf
[Service]
ExecStartPost=/usr/bin/systemctl daemon-reload
